# rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
source /home/ec2-user/CNV_RenderFarm/BIT/java-config.sh
export CLASSPATH="/home/ec2-user/CNV_RenderFarm/BIT/:/home/ec2-user/CNV_RenderFarm/BIT/samples/:/home/ec2-user/CNV_RenderFarm/src/raytracerBIT/:/home/ec2-user/CNV_RenderFarm/src/:./:$CLASSPATH"
cd /home/ec2-user/CNV_RenderFarm
java raytracer.WebServer
